# DEVELOPMENT/python3/POSTLOAD/usr/etc/profile.d/python3.sh
pathappend %PREFIX/bin			PATH
export PYTHON3=%PREFIX/bin/python3
