# DEVELOPMENT/python3/POSTLOAD/usr/etc/profile.d/python3.sh
pathappend %PREFIX/bin			PATH
pathappend %PREFIX/share/man		MANPATH
pathappend %PREFIX/lib/pkgconfig	PKG_CONFIG_PATH

export PYTHON3=%PREFIX/bin/python3
