# DESKTOP/libcddb/POSTLOAD/usr/etc/profile.d/libcddb.sh
pathappend	%PREFIX/bin			PATH
pathappend	%PREFIX/lib/pgkconfig		PKG_CONFIG_PATH



