# ICONS/adwaita-icon-theme/POSTLOAD/usr/etc/profile.d/adwaita-icon-theme.sh
pathappend	%PREFIX/share/pkgconfig	PKG_CONFIG_PATH
