# DESKTOP/figlet/POSTLOAD/usr/etc/profile.d/figlet.sh
pathappend	%PREFIX/bin 		PATH
pathappend	%PREFIX/share/man	MANPATH


