# CORE/isl/POSTLOAD/usr/etc/profile.d/isl.sh
pathappend %PREFIX/isl/include  	C_INCLUDE_PATH
pathappend %PREFIX/isl/include		CPLUS_INCLUDE_PATH
pathappend %PREFIX/isl/lib		LIBRARY_PATH
