# DESKTOP/potrace/POSTLOAD/usr/etc/profile.d/potrace.sh
pathappend	%PREFIX/bin		PATH
pathappend	%PREFIX/share/man	MANPATH
pathappend	%PREFIX			CMAKE_PREFIX_PATH

