# DEVELOPMENT/cmake/POSTLOAD/usr/etc/profile.d/cmake.sh
pathappend	%PREFIX/bin		PATH
pathappend	%PREFIX/share/aclocal	ACLOCAL_PATH