# CORE/gnutls/POSTLOAD/usr/etc/profile.d/gnutls.sh
pathappend	%PREFIX/lib		LIBRARY_PATH
pathappend	%PREFIX/include		C_INCLUDE_PATH
pathappend 	%PREFIX/bin		PATH

