# DESKTOP/gdbm/POSTLOAD/usr/etc/profile.d/gdbm.sh
pathappend	%PREFIX/bin		PATH
pathappend	%PREFIX/include		CPATH
pathappend	%PREFIX/lib		LIBRARY_PATH
pathappend	%PREFIX/share/man	MANPATH


