# CORE/bzip2/POSTLOAD/usr/etc/profile.d/bzip2.sh
pathappend %PREFIX/bin 		PATH
pathappend %PREFIX/include	C_INCLUDE_PATH
pathappend %PREFIX/include	CPLUS_INCLUDE_PATH
pathappend %PREFIX/lib		LIBRARY_PATH