# CORE/ca-certificates/POSTLOAD/usr/etc/profile.d/ca-certificates.sh
pathappend	%PREFIX/bin	PATH

export LMBX_CA_CERT=/usr/share/ca-cert/ca-bundle.crt
