# CORE/lvm2/POSTLOAD/usr/etc/profile.d/lvm2.sh
pathappend	%PREFIX/bin 				PATH
pathappend	%PREFIX/lib/python3.6/site-packages	PYTHONPATH
pathappend	%PREFIX/lib/pkgconfig			PKG_CONFIG_PATH
pathappend	%PREFIX/share/man			MANPATH
