# CORE/util-linux/POSTLOAD/usr/etc/profile.d/util-linux.sh
pathappend %PREFIX/lib		LIBRARY_PATH
pathappend %PREFIX/include	C_INCLUDE_PATH
pathappend %PREFIX/bin		PATH
pathappend %PREFIX/share/man	MANPATH
pathappend %PREFIX/lib/pkgconfig PKG_CONFIG_PATH



