# CORE/ncurses/POSTLOAD/usr/etc/profile.d/ncurses.sh
pathappend	%PREFIX/lib			LIBRARY_PATH
pathappend	%PREFIX/include		C_INCLUDE_PATH
pathappend	%PREFIX/include/ncursesw	C_INCLUDE_PATH
pathappend	%PREFIX/bin			PATH

