# CORE/glibc/POSTLOAD/usr/etc/profile.d/glibc.sh
pathappend 	%PREFIX/bin	PATH
pathappend	%PREFIX/lib	LIBRARY_PATH
