# DESKTOP/jasper/POSTLOAD/usr/etc/profile.d/jasper.sh
pathappend	/usr/pkg/jasper/bin		PATH
pathappend	/usr/pkg/jasper/lib/pkgconfig	PKG_CONFIG_PATH
