# CORE/bash/POSTLOAD/usr/etc/profile.d/bash.sh
pathappend	%PREFIX/bin		PATH
pathappend	%PREFIX/share/man	MANPATH
pathappend	%PREFIX/lib/pgkconfig	PKG_CONFIG_PATH
