# CORE/readline/POSTLOAD/usr/etc/profile.d/readline.sh
pathappend 	%PREFIX/include		C_INCLUDE_PATH
pathappend	%PREFIX/lib		LIBRARY_PATH
pathappend	%PREFIX/bin		PATH
pathappend	%PREFIX/share/man	MANPATH
