# CORE/python2/POSTLOAD/usr/etc/profile.d/python2.sh
pathappend %PREFIX/bin			PATH
export PYTHON2=%PREFIX/bin/python2
